# templates/default/service/init.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          faraday
# Required-Start:    $remote_fs $local_fs $network
# Required-Stop:     $remote_fs $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts and daemonize Faraday server
# Description:       Starts and daemonize Faraday server
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="Faraday server"
NAME=faraday
USER=$NAME
DAEMON="/opt/$NAME/server"
PIDFILE="/run/$NAME/$NAME.pid"
DAEMON_ARGS="--gui=no-gui"

[ -r /etc/default/$NAME ] && . /etc/default/$NAME

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

# Source function library.
<% case node['platform_family'] -%>
<% when 'rhel', 'fedora' -%>
. /etc/rc.d/init.d/functions
<% else -%>
. /lib/lsb/init-functions
<% end -%>

if [ ! -d "$(dirname $PIDFILE)" ]; then
    mkdir -p "$(dirname $PIDFILE)"
    chown $USER:$USER "$(dirname $PIDFILE)"
    chmod 755 "$(dirname $PIDFILE)"
fi

do_start()
{
    log_daemon_msg "Starting $DESC" "$NAME "
    [ -e $PIDFILE ] && PID=$(cat "$PIDFILE")
    if ( [ -e $PIDFILE ] && ps -p $PID 1>&2 > /dev/null )
    then
        log_action_end_msg 1 "already running, PID's $PID"
        exit 0
    elif ( [ -w $PIDFILE ] )
    then
        log_warning_msg "PID file found while ${NAME} is not running, removing file."
        rm $PIDFILE
    fi

    if [ "$RUN" != "true" ]; then
        log_action_msg "Not starting faraday: disabled by /etc/default/$NAME".
        exit 0
    fi

    start-stop-daemon --start --chuid $USER --pidfile=$PIDFILE --quiet -m -b --exec $DAEMON -- $DAEMON_ARGS || return 2
    [ -e "$PIDFILE" ] && chown $USER $PIDFILE
    return 0
}

do_stop()
{
    log_daemon_msg "Stopping $DESC" "$NAME "
    if [ ! -w $PIDFILE ]
    then
        log_warning_msg "PID file not found"
        return 4
    fi
    start-stop-daemon --stop --oknodo --chuid $USER --pidfile=$PID --user=$USER --pidfile="$PIDFILE" --retry=30 || return 1
    /bin/rm -f $PIDFILE
    return 0
}

case "$1" in
  start)
    do_start
    case "$?" in
        0|1) log_end_msg 0 ;;
        2) log_end_msg 1 ;;
    esac
    ;;
  stop)
    do_stop
    case "$?" in
        0) log_end_msg 0 ;;
        1) log_end_msg 1 ;;
    esac
    ;;
  status)
    status_of_proc -p "$PIDFILE" "$DAEMON" "$NAME"
    ;;
  restart|force-reload)
    do_stop
    case "$?" in
      0)
        log_end_msg 0
        do_start
        case "$?" in
            0) log_end_msg 0 ;;
            *) log_end_msg 1 ;;
        esac
        ;;
      *)
        log_end_msg 1
        ;;
    esac
    ;;
  *)
    echo "Usage: invoke-rc.d $NAME {start|stop|status|restart|force-reload}" >&2
    exit 3
    ;;
esac
